# claude-assistant.service
[Unit]
Description=Claude Office Assistant — Flask API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/claude-office-assistant/backend
EnvironmentFile=/var/www/claude-office-assistant/config/.env
ExecStart=/var/www/claude-office-assistant/venv/bin/gunicorn \
    app:app \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --timeout 120 \
    --access-logfile /var/log/claude-assistant/access.log \
    --error-logfile  /var/log/claude-assistant/error.log \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-assistant

[Install]
WantedBy=multi-user.target
